# sv/cloud-api-service-db/finish
#!/usr/bin/env bash
set -euo pipefail

source "$(dirname "$0")/../../env"
cd "${apps_base}/cloud-api-service"
eval "$(grep export .envrc)"
make db-stop 2>&1
