# sv/registry/run
#!/usr/bin/env bash
set -euo pipefail

source "$(dirname "$0")/../../env"
eval "$(grep export .envrc)"
npm run dev 2>&1
